# save-backlight
#!/sbin/openrc-run

# Copyright (C) 2019 Mattéo Rossillol‑‑Laruelle <beatussum@protonmail.com>
#
# This program is free software:  you can redistribute it and/or modify
# it under the terms of the GNU General Public License as published  by
# the Free Software Foundation, either version 3 of the License, or (at
# your option) any later version.
#
# This program is distributed in the hope that it will be useful,
# but WITHOUT ANY WARRANTY;  without even the implied warranty of
# MERCHANTABILITY or FITNESS FOR A PARTICULAR  PURPOSE.  See  the
# GNU General Public License for more details.
#
# You should have received a copy of the GNU General Public License along
# with this program. If not, see <https://www.gnu.org/licenses/>.


description="An OpenRC service for restoring the last brightness level on reboot"

cache_dir="/var/cache/save-backlight"

start() {
	ebegin "Restoring the screen brightness"

	checkpath -D -o "root:root" "${cache_dir}"

	for card in $(basename /sys/class/backlight/*); do
		if [ -e "${cache_dir}/${card}" ]; then
			local -r blight=$(cat "${cache_dir}/${card}")
		fi

		if [ ${blight:=0} -ne 0 ]; then
			echo ${blight} > "/sys/class/backlight/${card}/brightness"
		fi

		sts=$?
		test ${sts} -eq 0 || break
	done

	eend ${sts} "Cannot restore the brightness of ${card}"
}

stop() {
	ebegin "Saving the screen brightness"

	checkpath -D -o "root:root" "${cache_dir}"

	for card in $(basename /sys/class/backlight/*); do
		cp "/sys/class/backlight/${card}/brightness" \
			"${cache_dir}/${card}"

		sts=$?
		test ${sts} -eq 0 || break
	done

	eend ${sts} "Cannot save the brightness of ${card}"
}
